# proxmox/systemd/services/backup-vms-and-lxcs.service
[Unit]
Description=Performs backup of specified VMs and LXCs
After=zfs-mount.service

[Service]
Type=oneshot
ExecStart=create-vm-lxc-backups.sh