# _obsolete/systemd/services/power-outage-test.service
[Unit]
Description=Monitor power status via network beacons
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/power-outage-test.sh
User=root
